# deployments/systemd/go-wx.service
[Unit]
# Description of the service
Description=Go Weather Station Service
# Wait for network to be available before starting
After=network.target

[Service]
# Run as a specific user and group (change as needed)
User=go-wx
Group=go-wx

# Working directory
WorkingDirectory=/opt/go-wx

# Command to execute
ExecStart=/opt/go-wx/bin/go-wx -config /opt/go-wx/config/config.yaml

# Restart policy
Restart=always
RestartSec=5

# Environment variables
Environment=TZ=UTC

# Logging
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true

[Install]
# Start at boot
WantedBy=multi-user.target 